# edgeclip.service
# Systemd unit for EdgeClip (host-based Paperclip + Docker infrastructure).
#
# Install:
#   sudo cp edgeclip.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now edgeclip
#
# Logs:
#   journalctl -u edgeclip -f
#   tail -f /root/edgeclip/paperclip.log
#
# Adjust WorkingDirectory and User if your layout differs.

[Unit]
Description=EdgeClip — Paperclip Agent Orchestration
After=network-online.target docker.service
Requires=docker.service
Wants=network-online.target

[Service]
Type=forking
PIDFile=/root/edgeclip/.paperclip.pid
User=root
WorkingDirectory=/root/edgeclip
ExecStart=/root/edgeclip/start-service.sh
ExecStop=/root/edgeclip/stop.sh
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
